# vendor/code.uber.internal/engsec/wonka-go.git/debian/wonkad.postinst
#!/bin/sh

case "$1" in
    configure)
        systemctl enable wonkad.service || true
        service wonkad start || true
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        exit 0
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac
