# vendor/code.uber.internal/engsec/wonka-go.git/debian/wonkad.prerm
#!/bin/sh

case "$1" in
    remove|upgrade|deconfigure)
        systemctl disable wonkad.service || true
        service wonkad stop || true
        ;;
    failed-upgrade)
        echo "failed upgrade" || true
        ;;
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac
