# vendor/code.uber.internal/engsec/wonka-go.git/wonkad.service
[Unit]
Description=wonkad
After=syslog.target network.target

[Service]
User=root
Type=simple
ExecStart=/usr/sbin/wonkad
TimeoutStartSec=10s
TimeoutStopSec=10s
RestartSec=1s
Restart=on-failure
SuccessExitStatus=0 2
LimitNOFILE=65536
SyslogIdentifier=wonkad
CPUQuota=100%
MemoryLimit=4G

[Install]
WantedBy=multi-user.target
